# deploy/systemd/refund-python.service
[Unit]
Description=Refund Decision Agent Python FastAPI service
After=network-online.target refund-java.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/refund-decision-agent/app/python-agent
EnvironmentFile=/opt/refund-decision-agent/env/refund-python.env
ExecStart=/home/ubuntu/miniconda3/envs/refund-agent/bin/uvicorn app.main:app --host 0.0.0.0 --port 8001
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
